# infra/systemd/chimera-worker@.service
[Unit]
Description=Chimera RQ worker #%i
Documentation=file:/opt/chimera/README.md
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=chimera
Group=chimera
WorkingDirectory=/opt/chimera
EnvironmentFile=/etc/chimera/chimera.env
ExecStart=/opt/chimera/.venv/bin/rq worker --url ${REDIS_URL} --name chimera-worker-%i high normal low
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/chimera/logs/worker-%i.stdout.log
StandardError=append:/opt/chimera/logs/worker-%i.stderr.log

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chimera/logs /opt/chimera/storage
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
